# watchdog_mbim
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99

PROG="/root/bin/watchdog_mbim"
PID_FILE="/var/run/watchdog_mbim.pid"
CONFIG_FILE="/etc/watchdog_mbim.conf"
LOGHEADER="watchdog_mbim"

# Optional debug files (enable via touching)
WATCHDOG_DEBUG_FILE="/tmp/watchdog/mbim_debug"

start_service() {
    # Required binary check
    if [ ! -x "$PROG" ]; then
        logger -t "$LOGHEADER" "Cannot start: executable not found or not executable ($PROG)"
        return 1
    fi

    # Verify that the configuration file exists
    if [ ! -f "$CONFIG_FILE" ]; then
        logger -t "$LOGHEADER" "Cannot start: configuration file not found ($CONFIG_FILE)"
        return 1
    fi

    logger -t "$LOGHEADER" "Starting watchdog_mbim service"

    # Ensure runtime directory exists.
    # The daemon also creates it internally, but doing it here guarantees
    # the debug file can be touched before the process starts.
    mkdir -p /tmp/watchdog 2>/dev/null

    # Clean up stale files and lock directories from previous crash/force kill.
    # Lock dirs (/tmp/lock/watchdog_mbim.lock, /tmp/watchdog/ping_lock) are created
    # with atomic mkdir - use rmdir (not rm -rf) to avoid accidentally deleting a
    # non-empty directory if the path is ever misconfigured.
    rmdir /tmp/lock/watchdog_mbim.lock 2>/dev/null
    rmdir /tmp/watchdog/ping_lock 2>/dev/null
    rm -f /tmp/watchdog/ping_failed 2>/dev/null
    rm -f /tmp/watchdog/ping_once_force_fail 2>/dev/null
    rm -f /tmp/watchdog/atcmd_block 2>/dev/null
    rm -f /tmp/watchdog/ifup_lock 2>/dev/null
    rm -f /tmp/watchdog/ifup_lock_date 2>/dev/null

    # Enable debug logging for watchdog_mbim.
    # The presence of this file enables verbose syslog output from the daemon.
    # Uncomment the following lines to enable verbose debugging.
    #touch "$WATCHDOG_DEBUG_FILE"
    #logger -t "$LOGHEADER" "watchdog_mbim debug enabled ($WATCHDOG_DEBUG_FILE)"

    procd_open_instance
    procd_set_param command "$PROG"
    procd_set_param respawn 3 10 60
    procd_set_param pidfile "$PID_FILE"
    procd_set_param term_timeout 60
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

stop_service() {
    logger -t "$LOGHEADER" "Stopping watchdog_mbim service"
}

help() {
    echo
    echo "WATCHDOG_MBIM - minimal init.d wrapper"
    echo "-------------------------------------"
    echo "Commands:"
    echo "  start    Start the watchdog daemon"
    echo "  stop     Stop the watchdog daemon"
    echo "  restart  Restart the watchdog daemon"
    echo "  enable   Enable autostart at boot"
    echo "  disable  Disable autostart at boot"
    echo "  help     Show this message"
    echo
    echo "Configuration file: $CONFIG_FILE"
    echo
    echo "Debug mode: $( [ -f "$WATCHDOG_DEBUG_FILE" ] && echo ENABLED || echo DISABLED )"
    echo
    echo "To enable debug logging:"
    echo "  touch $WATCHDOG_DEBUG_FILE"
    echo
    echo "To disable debug logging:"
    echo "  rm -f $WATCHDOG_DEBUG_FILE"
    echo
}

# Reference:
# https://forum.openwrt.org/t/procd-init-interface-trigger-for-non-daemon-executables-that-are-not-running-permanently/58258/8
#
# Example procd/netifd interface trigger kept for reference only.
# Not enabled for watchdog_mbim because LTE/MBIM connectivity does not
# reliably map to netifd interface state events on all platforms.
#
# Test:
#   ifdown wan
#   ubus monitor | grep '"method":"up","data":{"interface":"wan"}}'
#
#service_triggers() {
#    procd_open_trigger
#    procd_add_interface_trigger "interface.*.up" "wan" "/etc/init.d/watchdog_mbim restart"
#    procd_close_trigger
#}
